# vineyard_robot.service
[Unit]
Description=Vineyard Robot Controller
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Working Code/Updated FULL Code Standalone
ExecStart=/usr/bin/python3 "/home/pi/Working Code/Updated FULL Code Standalone/main.py"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
